# ops/release-mirror/ecom-monitor-release-sync.timer
[Unit]
Description=Check for new e-commerce monitor releases

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true
Unit=ecom-monitor-release-sync.service

[Install]
WantedBy=timers.target
